# package/gluon-setup-mode/files/lib/gluon/setup-mode/rc.d/S50ffgt-tell-setup
#!/bin/sh /etc/rc.common

START=14

/lib/gluon/config-mode/notify-setup.sh
